# rootfs/etc/s6-overlay/s6-rc.d/unoserver-rest-api/run
#!/command/with-contenv sh
unoserver-rest-api --unoconvert-bin xunoconvert
